# 1.6/config/supervisord.conf
[supervisord]
nodaemon=true

[program:pgbouncer]
command=/bin/bash -c "/usr/local/bin/confd -onetime -backend env && /usr/local/bin/pgbouncer /etc/pgbouncer/pgbouncer.ini"
user=postgres